# templates/tomcat.service
[Service]
ReadWritePaths={{ tomcat_webapps }}

