# templates/tomcat.service
[Service]
ReadWritePaths={{ tomcat_webapps }}
ReadWritePaths=/srv/data
